# deploy/agent-memory-curate.service
[Unit]
Description=Multi-Agent Memory — Nightly Curation
After=network.target

[Service]
Type=oneshot
User=tgds
WorkingDirectory=/opt/multi-agent-memory
EnvironmentFile=/opt/multi-agent-memory/.env
Environment=PYTHONPATH=/opt/multi-agent-memory
ExecStart=/opt/multi-agent-memory/.venv/bin/python scripts/curate.py --limit 50
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
